# errata-rails/config/crontab
0 */3 * * * /usr/kerberos/bin/kinit -k errata/beehive\@REDHAT.COM -t /etc/errata/errata.keytab 2>&1
