# errata-rails/docker/12-jenkins_dind_slave/mkloops.service
[Unit]
Description="Make loop devices for docker-in-docker"
ConditionPathExists=/dev/loop-control

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/mkloops
